# etc/cron.d/shutdowncheck
*/10 * * * * root python3 /root/shutdowncheck.py --shutdowncheck
